# contest-reminders.timer
[Unit]
Description=Run contest reminder sync scripts daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=contest-reminders.service

[Install]
WantedBy=timers.target
